# triton-openbsd-guesttools/etc/rc.local
#!/usr/local/bin/bash
#
# Copyright (c) 2017 Joyent Inc., All rights reserved.

# Log output
LOG='/var/log/triton.log'
if [[ ! -e $LOG ]] ; then
  touch $LOG
fi

(/lib/smartdc/set-hostname)
(/lib/smartdc/set-root-authorized-keys)
if [[ ! -f /lib/smartdc/.firstboot-complete-do-not-delete ]] ; then
  (/lib/smartdc/firstboot)
fi
(/lib/smartdc/add-network-interface)
(/lib/smartdc/run-operator-script)
(/lib/smartdc/get-user-data)
(/lib/smartdc/run-user-script)
